# tar/config/.config/autostart/fix-x11-input.desktop
[Desktop Entry]
Type=Application
Name=Fix X11 Input
# Termux:X11 입력 초기화 — 세션 시작 시 1회 실행
#
# 근본 원인 (GitHub termux/termux-x11 #781):
#   Android는 다른 앱 전환 시 키 릴리스 이벤트를 보내지 않아
#   Alt 키가 X서버에 눌린 채 고착됨.
#   → Alt+클릭 = 우클릭, Alt+방향키 = 오작동
#
# 영구 해결: Termux:X11 설정 → Keyboard → "Capture DEX meta keys" 활성화
# 임시 우회: Super+I 단축키로 수동 리셋
Exec=bash -c 'sleep 5; xdotool key alt 2>/dev/null || true; xdotool key shift 2>/dev/null || true; xdotool key ctrl 2>/dev/null || true; xinput set-button-map "Virtual core pointer" 1 2 3 4 5 6 7 8 9 10 11 12 13 2>/dev/null || true; xset r rate 300 30 2>/dev/null || true; setxkbmap -option 2>/dev/null || true'
X-GNOME-Autostart-enabled=true
NoDisplay=true
StartupNotify=false
Terminal=false
